# cookbooks/logstash-server/templates/default/upstart-redis-server.conf.erb
# This file is managed by Chef
# Any changes made by hand will be overwritten

# Note, you need a 'redis' user on your machine to run this script.
# redis-server upstart init script
# this should log to /var/log/upstart/logstash.log
# elasticsearch_server: 127.0.0.1:9200

description     "redis-server"

setuid redis
setgid redis

start on virtual-filesystems or runlevel [2345]
stop on runlevel [!2345]

# Respawn it if the process exits
respawn
respawn limit 5 30
limit nofile 65550 65550

exec /usr/bin/redis-server /etc/redis/redis.conf 